# deploy/web-reader.service
[Unit]
Description=Web Reader

After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=web-reader
Group=web-reader
WorkingDirectory=/opt/web-reader
EnvironmentFile=/etc/web-reader/web-reader.env
ExecStart=/opt/web-reader/web-reader
Restart=on-failure
RestartSec=3s
TimeoutStopSec=15s
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
